# srv/salt/runit/files/etc/service/sumologic/run
#! /bin/sh
exec 2>&1
source /etc/envvars
/usr/local/bin/sumologic
exec /opt/SumoCollector/collector console -- -t  > /var/log/collector.log
