# srv/salt/harden/files/etc/cron.daily/aide.cron
05 4 * * * root /usr/sbin/aide --check > /var/log/aide.log 2>&1
