# srv/salt/runit/files/etc/service/apache/run
#! /bin/sh
# apache default; we typically will override this with out application specific stuff
exec 2>&1
/usr/sbin/apache2 -D FOREGROUND
