# srv/salt/runit/files/etc/service/logrotate/run
#! /bin/sh
exec 2>&1
source /etc/envvars
exec /usr/sbin/crond -n -s
